# deploy_tools/gunicorn-systemd.template.service
[Unit]
Description=Gunicorn server for stagingottg.live

[Service]
Restart=on-failure
User=ansuresh
WorkingDirectory=/home/ansuresh/sites/stagingottg.live
EnvironmentFile=/home/ansuresh/sites/stagingottg.live/.env

ExecStart=/home/ansuresh/sites/stagingottg.live/virtualenv/bin/gunicorn \
    --bind unix:/tmp/stagingottg.live.socket \
    djangonewproject.wsgi:application

[Install]
WantedBy=multi-user.target
